# supervisord.conf
[supervisord]
nodaemon=true

[program:ss-server]
command=python server.py -c /config.json > /dev/null

[program:kcptun]
command=/opt/kcptun/server_linux_amd64 %(ENV_KCPTUN_ARGS)s
